# server/txform-provisioner.timer
# ============================================================
# Txform.ph — systemd timer for the provisioner
#
# Runs the oneshot provisioner every 2 minutes, so access changes made in
# the portal reach Books quickly. The portal tells owners "usually under
# two minutes" — if you change OnUnitActiveSec, change that copy too.
#
#   sudo systemctl enable --now txform-provisioner.timer
#   systemctl list-timers txform-provisioner.timer     # confirm it ticks
#   journalctl -u txform-provisioner.service -n 50     # what it did
# ============================================================
[Unit]
Description=Run the Txform.ph provisioner periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
# A missed window (box asleep, deploy mid-flight) runs once on the next
# start rather than being skipped — a queued revoke must not be lost.
Persistent=true
Unit=txform-provisioner.service

[Install]
WantedBy=timers.target
